# infra/systemd/caiji-mvp.service
[Unit]
Description=caiji-mvp - QQ 群消息采集 FastAPI 服务
Documentation=https://github.com/l6444585-afk/caiji-mvp
After=network.target

[Service]
Type=exec
User=root
WorkingDirectory=/opt/fengyun/caiji-mvp
Environment="PATH=/opt/fengyun/caiji-mvp/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=/etc/caiji-mvp.env
ExecStart=/opt/fengyun/caiji-mvp/.venv/bin/python -m uvicorn server:app --host 127.0.0.1 --port 8090
Restart=always
RestartSec=3
StandardOutput=append:/var/log/caiji-mvp.log
StandardError=append:/var/log/caiji-mvp.log

# 安全加固
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
